# src/debpkg/DEBIAN/postinst
#!/bin/sh
set -e
chmod 544 /usr/bin/vcool
chmod 644 /etc/vcool/vcool.stg
chmod 755 /etc/init.d/vcool
update-rc.d vcool defaults
exit 0
